# package/etc/cron.weekly/plus-us-update
#!/bin/bash

apt-get --yes --allow-unauthenticated -q update;
apt-get --yes --allow-unauthenticated -q upgrade;
